# pkgs/udev-custom-callback/99-custom-bluetooth.rules
ACTION=="add", SUBSYSTEM=="input", ATTR{name}=="Jabra Elite 8 Active (AVRCP)", \
RUN+="udev-custom.sh bluetooth %E{NAME}"
